# deploy/com.sahilai.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sahilai.worker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/opt/sahil-ai/run_worker.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/opt/sahil-ai</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/sahil-worker.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/sahil-worker.log</string>
</dict>
</plist>
